# scripts/prolog-security-backup.service
[Unit]
Description=Backup Prolog Security Expert System

[Service]
Type=oneshot
User=micu
Group=micu
WorkingDirectory=/home/micu/prolog
ExecStart=/home/micu/prolog/scripts/backup_project.sh
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/micu/prolog/data/backups
