# scripts/prolog-security-backup.service
[Unit]
Description=Backup Prolog Security Expert System

[Service]
Type=oneshot
User=micu
Group=micu
WorkingDirectory=/home/micu/prolog
ExecStart=/home/micu/prolog/scripts/backup_project.sh
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=read-only
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX
ProtectProc=invisible
ProcSubset=pid
ReadWritePaths=/home/micu/prolog/data/backups
